# playbooks/redis/files/redis-6380.service
[Unit]
Description=Redis (6380)
Documentation=https://redis.io/docs/
After=syslog.target

[Service]
Type=simple
User=root
Group=root
ExecStart=/opt/redis/bin/redis-server /opt/redis/etc/redis-6380.conf
KillSignal=15
SuccessExitStatus=143

[Install]
WantedBy=multi-user.target
